# systemd/hermes-tasklane-sync.service
[Unit]
Description=hermes-tasklane sync
After=default.target

[Service]
Type=oneshot
ExecStart=__EXECUTABLE__ --config __CONFIG_PATH__ sync
