# systemd/hermes-tasklane-dashboard.service
[Unit]
Description=hermes-tasklane dashboard
After=default.target

[Service]
Type=simple
ExecStart=__EXECUTABLE__ --config __CONFIG_PATH__ dashboard --host __DASHBOARD_HOST__ --port __DASHBOARD_PORT__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
