# systemd/hermes-tasklane-watch.service
[Unit]
Description=hermes-tasklane watch
After=default.target

[Service]
Type=oneshot
ExecStart=__EXECUTABLE__ --config __CONFIG_PATH__ watch --quiet-ok
