# systemd/hermes-tasklane-reconcile.service
[Unit]
Description=hermes-tasklane reconcile
After=default.target

[Service]
Type=oneshot
ExecStart=__EXECUTABLE__ --config __CONFIG_PATH__ reconcile
